# docker/app/supervisord.d/gunicorn.ini
[program:gunicorn]
command=gunicorn --config /etc/gunicorn/config.py mysite.wsgi:application
startretries=1000000 # only for development
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
